# services/farm3d.service.in
[Unit]
Description=farm3d
After=moonraker.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User={{PRINTER_USER}}
WorkingDirectory={{PRINTER_HOME}}/farm3d/
ExecStart={{PRINTER_HOME}}/farm3d/run.sh

[Install]
WantedBy=multi-user.target
